# container/overlay/etc/systemd/system/console-getty.service.d/autologin.conf
# Autologin the player on the nspawn console (/dev/console). The blank
# ExecStart= is required first: it clears the unit's original line so the
# second ExecStart fully replaces it.
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin player --noclear - $TERM
